# build_dir/target-mipsel_24kc_musl/root.orig-ramips/usr/lib/opkg/info/https-dns-proxy.list
/etc/config/https-dns-proxy
/etc/init.d/https-dns-proxy
/usr/sbin/https-dns-proxy
