# build_dir/target-mipsel_24kc_musl/root.orig-ramips/usr/lib/opkg/info/watchcat.list
/etc/uci-defaults/50-watchcat
/lib/upgrade/keep.d/watchcat
/usr/bin/watchcat.sh
/etc/init.d/watchcat
